# wiki/meduseld-wiki-scrape.service
[Unit]
Description=Meduseld Wiki Scrape (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=vertebra
Environment="WIKI_DIR=/srv/wiki/icarus"
ExecStart=/srv/apps/meduseld/venv/bin/python /srv/apps/meduseld/wiki/scrape_wiki.py
TimeoutStartSec=3600
